# deploy/opsagent.service
[Unit]
Description=Opsagent deployment management server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opsagent
Group=opsagent

ExecStart=/var/lib/opsagent/bin/opsagent

# Environment — override via /etc/opsagent/env or systemctl edit
EnvironmentFile=-/etc/opsagent/env
Environment=PATH=/nix/var/nix/profiles/default/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

Restart=on-failure
RestartSec=5

# Hardening
# NoNewPrivileges is intentionally not set: the runner uses sudo to
# restart its own systemd unit via the sudoers rule installed by the
# deploy script. This requires privilege escalation to succeed.
ProtectSystem=strict
ProtectHome=yes
# Allow writes to the data dir and the nix store
ReadWritePaths=/var/lib/opsagent
ReadWritePaths=/nix
# Allow writes to ACME cert cache
ReadWritePaths=/var/lib/opsagent/.certs
# Allow sudo to create its timestamp directory (needed for systemctl restart)
ReadWritePaths=/run/sudo
# Allow binding to privileged ports (80/443)
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID

[Install]
WantedBy=multi-user.target
